# debian/prerm
#! /bin/sh

set -e

case "$1" in
    remove|upgrade|deconfigure)
        systemctl stop fvapp-sabnzbd-rn.service
        systemctl disable fvapp-sabnzbd-rn.service || true
        ;;
    failed-upgrade)
        ;;
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
